# elasticluster/share/playbooks/roles/spark-yarn/files/etc/profile.d/spark-yarn.sh
# Additional settings for PySpark running over YARN

# force PySpark to use the Hadoop/YARN cluster
# (default is to use a local executor)
export MASTER=yarn-client
